# gnome/ctrl-rightclick/ctrl-rightclick.service
[Unit]
Description=ctrl-rightclick — hold Control + trackpad press = right-click (macOS-style)
Documentation=file://%h/code/chesterbr/dotfiles/gnome/ctrl-rightclick/README.md
# Needs Toshy's virtual keyboard to exist; the script also waits for it, and
# restarts cover Toshy (re)starting after us.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/code/chesterbr/dotfiles/gnome/ctrl-rightclick/ctrl-rightclick.py
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
